# board/raspberrypi/automount.rules
ACTION=="add",KERNEL=="sda2", RUN+="/usr/bin/ntfs-3g /dev/sda2 /ntfs"
ACTION=="remove",KERNEL=="sda2", RUN+="/bin/umount /ntfs"
